# bctmark/ansible/roles/ethereum/templates/parity.service.j2
[Unit]
Description=Ethereum Parity client
After=network.target

[Service]
User={{parity_user}}
Group={{parity_group}}

Type=simple
ExecStart=/usr/bin/parity --config /home/{{parity_user}}/config.toml
Restart=always
PrivateTmp=true
TimeoutStopSec=120s
TimeoutStartSec=20s
RestartSec=15s
StartLimitInterval=4h
StartLimitBurst=500

[Install]
WantedBy=multi-user.target